# supervisord.conf
[supervisord]
nodaemon=true

[program:rails]
command=bundle exec bin/rails s -e production
directory=/usr/src/rails
stdout_logfile=/usr/src/rails/log/%(program_name)s.log
stderr_logfile=/usr/src/rails/log/%(program_name)s.log
autorestart=true